# protected/data/container/webmysql/supervisord.conf
[supervisord]
nodaemon=true
[program:crond]
command = /usr/sbin/cron
user = root
autostart = true
[program:mysqld]
command = /usr/sbin/mysqld
user = root
autostart = true
